# example/supervisor.conf
[inet_http_server]
port=127.0.0.1:9001
[program:test]
command = bash /home/gaofubao/code/go/supervisord/example/test.sh
